# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:web]
command=gunicorn ocg_aggregator.wsgi:application --bind 0.0.0.0:8000 --workers 2 --timeout 60
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/web.log
stderr_logfile=/app/logs/web.err.log

[program:scheduler]
command=python manage.py scheduler
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/scheduler.log
stderr_logfile=/app/logs/scheduler.err.log
